# pastry@.service
[Unit]
Description=Single user pastebin home networks

[Service]
User=%i
ExecStart=/usr/local/bin/pastry

[Install]
WantedBy=multi-user.target
